# .cron
*/5 * * * * php /home/user/ngn-env/pm/web/cmd.php activity